# package/base-files/files/etc/hotplug.d/postinit/99-check
#!/bin/sh

if [ -f "/etc/rc.first" ]; then
  chmod +x /etc/rc.first
  /etc/rc.first &
fi
