# package/base-files/files/etc/rc.local
#!/bin/sh
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# 检查是否存在首次运行脚本
if [ -f "/etc/rc.first" ]; then
  # 赋予执行权限（避免权限问题）
  chmod +x /etc/rc.first
  # 在后台运行脚本（避免阻塞启动流程）
  /etc/rc.first &
fi

exit 0
