# package/base-files/files/etc/hotplug.d/iface/01-checkwan
#!/bin/sh

[ "$INTERFACE" = "wan" ] && [ "$ACTION" = "ifup" ] && sh /etc/check_lan_wan_ip.sh